# deploy/crontab
# ETL schedules
# Feed reader: daily at midnight
0 0 * * * /app/deploy/etl-run.sh feed_reader etl.handlers.feed_reader

# Catalog updater: daily at 03:00
0 3 * * * /app/deploy/etl-run.sh catalog_updater etl.handlers.catalog_updater

# Log rotation: daily at 04:00 — compress logs over 10 MB
0 4 * * * find /var/log/etl -name '*.log' -size +10M -exec sh -c 'mv "$1" "$1.$(date +\%Y\%m\%d).old" && gzip -f "$1.$(date +\%Y\%m\%d).old"' _ {} \;
